# .deploy/supervisord.conf
[supervisord]
nodaemon=true
user=root  ; root olarak çalıştırmayı bilinçli seçtiğimizizi bildir (uyarı susar)

[program:gunicorn]
; Dynamically sized & recycled workers for stability (max-requests to avoid memory bloat)
command=/bin/sh -c 'exec gunicorn -w ${GUNICORN_WORKERS:-2} \
	--max-requests ${GUNICORN_MAX_REQUESTS:-1000} \
	--max-requests-jitter ${GUNICORN_MAX_REQUESTS_JITTER:-50} \
	--timeout ${GUNICORN_TIMEOUT:-90} --bind 0.0.0.0:${GUNICORN_PORT:-8000} wsgi:app'
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=20

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=10
